# software/core/oqm-core-characteristics/build/installers/CharacteristicsDeb/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable "oqm-core-characteristics.service"
systemctl start "oqm-core-characteristics.service"
