# software/core/oqm-core-characteristics/build/installers/CharacteristicsDeb/DEBIAN/prerm
#!/bin/bash

systemctl disable "oqm-core-characteristics.service"
systemctl stop "oqm-core-characteristics.service"
